# supervisor/conf/jenkins.ini
[program:jenkins]
command=/bin/su jenkins /bin/bash -c "/usr/bin/java -Xms1024m -Xmx1024m -Xmsn1024m -server -Djava.awt.headless=true -jar /usr/local/lib/jenkins/jenkins.war --httpPort=8080"
process_name=%(program_name)s
pidfile=/opt/jenkins/system/process.pid
numprocs=1
autorestart=false
autostart=true
