# supervisor/conf/jenkins.ini
[program:jenkins]
command=/bin/su jenkins /bin/bash -c "JENKINS_HOME=/opt/jenkins/data JENKINS_ENABLE_ACCESS_LOG=yes /usr/bin/java -Xms1024m -Xmx1024m -server -Djava.awt.headless=true -jar /usr/local/lib/jenkins/jenkins.war --logfile=/opt/jenkins/log/jenkins.log --httpPort=8080"
process_name=%(program_name)s
pidfile=/opt/jenkins/system/process.pid
numprocs=1
autorestart=false
autostart=true
